# backend/config/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/blog/backend/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
